# market-monitor.service
[Unit]
Description=Market Monitor Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/macro-news
Environment=PYTHONPATH=/home/ubuntu/macro-news
ExecStart=/home/ubuntu/macro-news/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/market_monitor/output.log
StandardError=append:/var/log/market_monitor/error.log

[Install]
WantedBy=multi-user.target